# myapi.service
[Unit]
Description=My test api

[Service]
User=vagrant
ExecStart=/srv/test_api/main

[Install]
WantedBy=multi-user.target